# argocd-helm-secrets.service
[Unit]
Description="build argocd-helm-secrets"

[Service]
ExecStart=/home/user/src/argocd-helm-secrets/build.sh
